# agent/ansible/roles/noise_agent_restart/templates/agent-restart.service.j2
[Unit]
Description=Periodic assessment-agent restart (noise demo)

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart assessment-agent
