# agent/ansible/roles/noise_agent_restart/templates/agent-restart.timer.j2
[Unit]
Description=assessment-agent restart timer (noise demo)

[Timer]
OnBootSec=5min
OnUnitActiveSec={{ agent_restart_interval | default('5min') }}

[Install]
WantedBy=timers.target
