# agent/ansible/roles/noise_offline_once/templates/agent-offline.service.j2
[Unit]
Description=Stop assessment-agent once (noise_offline_once)

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop assessment-agent
