# agent/ansible/roles/noise_offline_once/templates/agent-offline.timer.j2
[Unit]
Description=Stop assessment-agent 5 min after boot (noise_offline_once)

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
